# .config/systemd/user/swaybg.service
[Unit]
PartOf=graphical-session.target
After=graphical-session.target
Requisite=graphical-session.target

[Service]
ExecStart=/usr/bin/swaybg -o "*" -m "fill" -i "%h/.config/backgrounds/wallpaper.jpg"
Restart=on-failure
